# rpi_ws281x_hub/rpi-ws281x-hub.service
[Unit]
Description=rpi-ws281x-hub-service
After=network.target

[Service]
ExecStart=rpi-ws281x-hub

[Install]
WantedBy=multi-user.target